# workspace/supervisord.d/docker.ini
[program:docker]
command=/usr/bin/dockerd --debug
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB