# workspace/supervisord.d/ttyd.ini
[program:ttyd]
command=/bin/sh -c "ttyd -p 8022 -O -P 15 /bin/bash"
user=ubuntu
directory=/home/ubuntu
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
